# deploy/systemd/myfreeceevee-parser.service
[Unit]
Description=MuhFweeCeeVee Parser (FastAPI)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/git-public/MuhFweeCeeVee-Professional CV Template Generator/services/parser
ExecStart=%h/git-public/MuhFweeCeeVee-Professional CV Template Generator/services/parser/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
